# script/processor
#!/sbin/openrc-run
command="node"
command_args_foreground="/code/processor.js > /var/log/processor.log 2>&1"
description="Hive wallet Processor Service"
pidfile="/var/run/processor.pid"
supervisor=supervise-daemon
