# script/server
#!/sbin/openrc-run
command="node"
command_args_foreground="/code/server.js > /var/log/server.log 2>&1"
description="Hive wallet Server"
pidfile="/var/run/server.pid"
supervisor=supervise-daemon
